# datasciencebox/salt/ipython/templates/notebook.conf
{%- from 'system/settings.sls' import user with context -%}

{% set environment = [] -%}

{%- if salt['cmd.retcode']('test -e /usr/lib/spark/conf/spark-env.sh') == 0 -%}
{%- from 'mesos/spark/settings.sls' import env with context -%}
{%- for var in env -%}
{%- do environment.append(var ~ '="' ~ env[var] ~ '"') -%}
{%- endfor -%}
{%- endif -%}

[program:notebook]
command=/home/{{ user }}/anaconda/bin/ipython notebook --ip=0.0.0.0 --port=8888 --no-browser
user={{ user }}
numprocs=1
directory=/home/{{ user }}/notebooks
autostart=false
autorestart=true
stdout_logfile=/var/log/ipython/notebook.log
stdout_logfile_maxbytes=50MB
redirect_stderr=true
environment={{ environment | join(',') }}
