# datasciencebox/salt/mesos/marathon/marathon.conf
{%- from 'java/settings.sls' import java with context -%}
{%- from 'cdh5/zookeeper/settings.sls' import zk with context -%}

[program:marathon]
command={{ java.java_home }}/bin/java -Djava.library.path=/usr/local/lib:/usr/lib:/usr/lib64 -Djava.util.logging.SimpleFormatter.format=%%2$s%%5$s%%6$s%%n -Xmx512m -cp /usr/bin/marathon mesosphere.marathon.Main --zk zk://{{ zk['connection_string'] }}/marathon --master zk://{{ zk['connection_string'] }}/mesos
startsecs=5
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/marathon.log
stdout_logfile_maxbytes=50MB
redirect_stderr=true
